# ods-devenv/crontab/crontab

55 */3 * * * /home/openshift/bin/restart_suite
